# now_playing.service
#+
#now_playing
#-

Description=Now Playing
After=network-online.target

[Service]
Type=simple
ExecStart=/root/now_playing/now_playing_forking.py --logger=debug
Restart=always
RestartSec=60
StartLimitIntervalSec=0

[Install]
WantedBy=multi-user.target
