# service/now_playing_feed.service
#+
#now_playing
#-

Description=Now Playing Fedd
After=network-online.target

[Service]
Type=simple
ExecStart=/root/now_playing/now_playing_feed.py --logger=debug
Restart=always
RestartSec=60
StartLimitIntervalSec=0

[Install]
WantedBy=multi-user.target
